# sway/wsl/win-to-clipboard.service
[Unit]
Description=Sync Windows clipboard to Wayland

[Service]
ExecStart=%h/.local/bin/win-to-clipboard
Restart=on-failure
RestartSec=3
